# .config/systemd/user/sway-swayidle.service
[Unit]

[Service]
ExecStart=/usr/bin/swayidle -w
Restart=on-failure
